# conf/rslsync
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

rsl_config="/etc/rslsync.conf"
rsl_command="/bin/rslsync"

depend() {
    need net
}

checkconfig() {
    if [ ! -f ${rsl_config} ]; then
        ewarn "${rsl_config} does not exist."
    fi

    return 0
}

start() {
    checkconfig || return 1

    ebegin "Starting Resilio Sync"
    su - pl -c "${rsl_command} --config ${rsl_config}"
    eend $?
}

stop() {
    ebegin "Stopping Resilio Sync"
    killall -q /bin/rslsync || retval=1
    #pkill --pidfile /var/run/resilio/resilio.pid
    eend $?
}


